# stands-28-dynamic_web/files/service/go.service
[Unit]
Description=Go

[Service]
Type=simple
Restart=always

ExecStart=/opt/websrv/websrv

[Install]
WantedBy=multi-user.target