# stands-28-dynamic_web/files/service/nodejs_starter.service
[Unit]
Description=Nodejs

[Service]
Type=simple
Restart=always

ExecStart=/bin/node /opt/static/webservice.js

[Install]
WantedBy=multi-user.target